# upm/ca/misc/supervisord.conf
[supervisord]
nodaemon=true

[program:bash]
command=bash
priority=100
startretries=0
