# upm/ca/misc/supervisord.conf
[supervisord]
nodaemon=true

[program:bash]
command=bash
priority=100
startretries=0

[program:default_gateway]
command=bash -c "eval $GW_REPLACE"
priority=1
startretries=0
