# org1/supervisord.conf
[supervisord]
nodaemon=true

[program:snmpd]
command=/usr/sbin/snmpd -f -u snmp -g snmp -Ismux,mteTrigger,mteTriggerConf -p /run/snmpd.pid
priority=100
