# ccalculo/supervisord.conf
[supervisord]
nodaemon=true

[program:snmpd]
command=/usr/sbin/snmpd -f -Lsd -Lf /dev/null -u snmp -g snmp -I -smux mteTrigger mteTriggerConf -p /run/snmpd.pid
priority=100
