# ccalculo/supervisord.conf
[supervisord]
nodaemon=true

[program:snmpd]
command=/etc/init.d/snmpd start
priority=100
startretries=0

[program:snmptrapd]
command=snmptrapd -f -Leo
priority=100
