# upm/voip/misc/supervisord.conf
[supervisord]
nodaemon=true

[program:snmpd]
command=/etc/init.d/snmpd start
priority=100
startretries=0

[program:asterisk]
command=asterisk -f -vvvvvvv
priority=100

