# um/db/misc/supervisord.conf
[supervisord]
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:postgres]
command=docker-entrypoint.sh postgres
priority=100

[program:firewall]
command=iptables-restore -n < /etc/iptables_firewall
priority=5
startretries=0

[program:default_gateway]
command=bash -c "eval $GW_REPLACE"
priority=1
startretries=0
